# xconfig/.config/systemd/user/xss-deactivate.timer
[Unit]
Description=Deactivate xss
RefuseManualStop=no

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=graphical.target