# xconfig/.config/systemd/user/xss-deactivate.service
[Unit]
Description=Deactivate xss
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/usr/bin/xscreensaver-command -deactivate